# dnsmasq.ini
[program:dnsmasq]
process_name = master
command=/usr/sbin/dnsmasq-entrypoint.sh
startsecs=0
autorestart=false