# inetd.ini
[program:inetd]
process_name = master
command=/usr/sbin/inetd-entrypoint.sh
startsecs=0
autorestart=false