# postfix.ini
[program:postfix]
process_name = master
command=/usr/sbin/postfix-entrypoint.sh
startsecs=0
autorestart=false